# server/deploy/wikilinker.service
# wikilinker/wikilinker.service
# Copy to /etc/systemd/system/wikilinker.service on the server

[Unit]
Description=Wikilinker Proxy Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/wikilinker
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=CACHE_DIR=/var/cache/wikilinker
Environment=CACHE_MAX_SIZE=1073741824
EnvironmentFile=-/etc/wikilinker/env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/cache/wikilinker
PrivateTmp=true

[Install]
WantedBy=multi-user.target
